# systemd/phpborg-workers.target
[Unit]
Description=phpBorg Worker Pool - All backup worker instances
After=phpborg-scheduler.service

[Install]
WantedBy=multi-user.target
