# systemd/phpborg-scheduler.service
[Unit]
Description=phpBorg Scheduler Worker - Handles backup schedules and periodic stats collection
After=network.target redis-server.service mariadb.service
Wants=redis-server.service mariadb.service

[Service]
Type=simple
User=phpborg
Group=phpborg
WorkingDirectory=/opt/newphpborg/phpBorg
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

# Main scheduler command
ExecStart=/usr/bin/php /opt/newphpborg/phpBorg/bin/phpborg scheduler:start

# Restart policy
Restart=always
RestartSec=5

# Resource limits
MemoryMax=512M
CPUQuota=50%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=phpborg-scheduler

# Security
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/newphpborg/phpBorg/var
ReadWritePaths=/opt/newphpborg/phpBorg/logs
ReadWritePaths=/opt/newphpborg/phpBorg/.git
ReadWritePaths=/opt/backups
ReadWritePaths=/var/lib/phpborg/backups
ReadWritePaths=/var/log/phpborg
ReadWritePaths=/tmp

[Install]
WantedBy=multi-user.target
