# systemd/phpborg-worker@.service
[Unit]
Description=phpBorg Worker #%i
Documentation=https://github.com/altzone/phpBorg
After=network.target mariadb.service redis.service
PartOf=phpborg-workers.target

[Service]
Type=simple
User=phpborg
Group=phpborg
WorkingDirectory=/opt/newphpborg/phpBorg

EnvironmentFile=/opt/newphpborg/phpBorg/.env
Environment="PHP_ENV=production"
Environment="WORKER_ID=%i"

ExecStart=/usr/bin/php8.3 /opt/newphpborg/phpBorg/bin/phpborg worker:start --queue=default --worker-id=%i

Restart=always
RestartSec=5s

StandardOutput=journal
StandardError=journal
SyslogIdentifier=phpborg-worker@%i

# No namespace isolation - FUSE mounts need global visibility
NoNewPrivileges=false

LimitNOFILE=65536
TimeoutStopSec=30s

[Install]
WantedBy=phpborg-workers.target
